# playbooks/roles/vexor/rx/templates/rc.local.j2
#!/bin/sh

{{ vx_home }}/bin/worker-rx >> /var/log/vxworker-rx.log 2>&1
